# nn.conf
[program:nn]
command=python3 /home/curator2/nn/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/nn/test.err.log
stdout_logfile=/home/curator2/nn/test.out.log
